# deploy/archive-transcriber.service
[Unit]
Description=Archive transcriber
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/livevtt
Environment=PYTHONUNBUFFERED=1
Environment=WORKERS=4
Environment=GPUS=0,1
ExecStart=/home/ubuntu/livevtt/scripts/run_archive_transcriber.sh
Restart=on-failure
RestartSec=30
KillMode=control-group
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
